# codemotion/nginx/files/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/etc/init.d/nginx start
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

[program:consul]
command=/usr/local/bin/consul agent -datacenter %(ENV_CONSUL_DATACENTER)s -config-dir /config/consul -client 0.0.0.0 -retry-join %(ENV_CONSUL_SERVER)s -bind 0.0.0.0 -grpc-port=8502 -data-dir /etc/consul -hcl "enable_central_service_config = true"
autostart=true
autorestart=true
stderr_logfile=/var/log/consul.err.log
stdout_logfile=/var/log/consul.out.log

[program:consul-template]
command=/usr/local/bin/consul-template -config=/config/consul_template/consul_template.hcl
autostart=true
autorestart=true
stderr_logfile=/var/log/consul-tempalte.err.log
stdout_logfile=/var/log/consul-template.out.log